# launchd/com.immich.thumbnail.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!-- EDIT PATHS BELOW: replace /path/to with your actual install location -->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.immich.thumbnail</string>
    <key>ProgramArguments</key>
    <array>
        <string>/path/to/immich-apple-silicon/ml/venv/bin/python3</string>
        <string>-m</string>
        <string>thumbnail</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/path/to/immich-apple-silicon</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>DB_HOST</key><string>localhost</string>
        <key>DB_PORT</key><string>5432</string>
        <key>DB_NAME</key><string>immich</string>
        <key>DB_USER</key><string>postgres</string>
        <key>DB_PASS</key><string>YOUR_DB_PASSWORD</string>
        <key>UPLOAD_DIR</key><string>/path/to/your/immich/upload</string>
        <key>PHOTOS_DIR</key><string>/path/to/your/photos</string>
        <!-- Uncomment if your Immich Docker uses non-standard volume mounts (e.g., Synology NAS):
        <key>CONTAINER_UPLOAD_PATH</key><string>/data/upload</string>
        <key>CONTAINER_PHOTOS_PATH</key><string>/mnt/media/Syno</string>
        -->
        <key>BATCH_SIZE</key><string>20</string>
        <key>POLL_INTERVAL</key><string>5</string>
    </dict>
    <key>KeepAlive</key><true/>
    <key>RunAtLoad</key><true/>
    <key>StandardOutPath</key><string>/tmp/immich-thumbnail.log</string>
    <key>StandardErrorPath</key><string>/tmp/immich-thumbnail.err</string>
</dict>
</plist>
