# launchd/com.immich.ffmpeg-proxy.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!-- EDIT PATHS BELOW: replace /path/to with your actual install location.
     Python must have Pillow installed (for WebP fallback).
     If using the immich-ml-metal venv, use its python3 path instead of /usr/bin/python3. -->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.immich.ffmpeg-proxy</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>/path/to/immich-apple-silicon/ffmpeg-proxy/server.py</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/path/to/immich-apple-silicon/ffmpeg-proxy</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>FFMPEG_PROXY_PORT</key>
        <string>3005</string>
        <key>UPLOAD_DIR</key>
        <string>/path/to/your/immich/upload</string>
        <key>PHOTOS_DIR</key>
        <string>/path/to/your/photos</string>
        <!-- Uncomment if your Immich Docker uses non-standard volume mounts (e.g., Synology NAS):
        <key>CONTAINER_UPLOAD_PATH</key>
        <string>/data/upload</string>
        <key>CONTAINER_PHOTOS_PATH</key>
        <string>/mnt/media/Syno</string>
        -->
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/bin:/bin</string>
    </dict>
    <key>KeepAlive</key>
    <true/>
    <key>RunAtLoad</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/tmp/immich-ffmpeg-proxy.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/immich-ffmpeg-proxy.err</string>
</dict>
</plist>
